# deploy/errordisplay.service
[Unit]
Description=Live Error Display Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/live-error-display
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=errordisplay

[Install]
WantedBy=multi-user.target
